# bin/systemd/celerybeat.service
# see https://docs.celeryproject.org/en/latest/userguide/daemonizing.html#daemon-systemd-generic

[Unit]
Description=Celery Beat Service
After=network.target

[Service]
Type=simple
User=live
Group=live
EnvironmentFile=/etc/conf.d/celerybeat.live.conf
WorkingDirectory=/home/live/edc
ExecStart=/bin/sh -c '${CELERY_BIN} beat  \
  -A ${CELERY_APP} --pidfile=${CELERYBEAT_PID_FILE} \
  --logfile=${CELERYBEAT_LOG_FILE} --loglevel=${CELERYD_LOG_LEVEL}'
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
